# deploy/celery-beat.service
[Unit]
Description=Travel Pro – Celery Beat Scheduler
After=network.target redis.service
Requires=redis.service

[Service]
Type=simple
User=travelpro
Group=travelpro
WorkingDirectory=/var/www/travelpro
Environment="DJANGO_SETTINGS_MODULE=config.settings.production"
EnvironmentFile=/var/www/travelpro/.env
ExecStart=/var/www/travelpro/.venv/bin/celery \
    -A config.celery beat \
    --pidfile=/run/travelpro/celery_beat.pid \
    --logfile=/var/log/travelpro/celery_beat.log \
    --loglevel=WARNING \
    --scheduler django_celery_beat.schedulers:DatabaseScheduler
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
